# packages/sysutils/systemd/system.d/openvpn.service
[Unit]
Description=OpenVPN Autorun Service

[Service]
Type=forking
Requires=network-online.service
After=network-online.service
ExecStart=/usr/sbin/openvpn --daemon --config /usr/share/openvpn/vpn.ovpn
Restart=always
RestartSec=15

[Install]
WantedBy=kodi.target
